# scripts/com.cinch.portal.monthlypull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.cinch.portal.monthlypull</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/homecomputer/Documents/Claude/Projects/Cinch Portal/sitelink-backend/scripts/monthly-pull.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key><integer>2</integer>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/homecomputer/Documents/Claude/Projects/Cinch Portal/sitelink-backend/logs/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/homecomputer/Documents/Claude/Projects/Cinch Portal/sitelink-backend/logs/launchd.err.log</string>
</dict>
</plist>
